# crontab
* * * * * root python3 /app/src/main.py
# Mandatory blank line